# config/deluge.conf
[program:deluge-script]
autostart=true
autorestart=true
startsecs = 0
user = abc
command = /scripts/deluge.sh
umask = 022
environment = PYTHON_EGG_CACHE="/config/.cache/Python-Eggs"
