# config/deluge.conf
[program:deluge-script]
autorestart = false
startsecs = 0
user = abc
command = /scripts/deluge.sh
umask = 022
